# home/.config/autostart/xhost_cups.desktop
[Desktop Entry]
Exec=/usr/bin/xhost +local:cups
Icon=
Name=xhost_cups
Path=
Terminal=False
Type=Application
